# configs/vagrant/teamtreats
#!/sbin/openrc-run
description="Team Treats webapp (TT)"
pidfile="/var/run/myapp.pid"

command="/usr/bin/dotnet"
command_args="/var/www/app/teamtreats-webapp.dll"
command_background="yes"
command_user="vagrant"
command_env="ASPNETCORE_ENVIRONMENT=Development"

depend() {
  need net
  use logger dns
}
